# daemon/data/earport-daemon.service
[Unit]
Description=EarPort AirPods Daemon
Documentation=https://github.com/Anoryth/earport
After=bluetooth.target

[Service]
Type=dbus
BusName=io.github.anoryth.EarPort
ExecStart=/usr/local/bin/earport-daemon
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# Creates $XDG_CONFIG_HOME/earport before sandboxing and keeps it writable
# (ReadWritePaths= would fail at namespace setup while the dir is missing)
ConfigurationDirectory=earport
PrivateTmp=true

[Install]
WantedBy=default.target
